# deploy/systemd/telegram-bot.service
[Unit]
Description=Telegram AI Bot Service
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/telegram-ai-system
Environment="PATH=/home/ubuntu/telegram-ai-system/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/home/ubuntu/telegram-ai-system"

# 使用虚拟环境中的 Python 运行 Telegram Bot
ExecStart=/home/ubuntu/telegram-ai-system/venv/bin/python main.py

# 如果不使用虚拟环境，使用系统 Python
# ExecStart=/usr/bin/python3 main.py

# 重启策略
Restart=always
RestartSec=10

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-bot

# 安全设置
NoNewPrivileges=true
PrivateTmp=true

# 资源限制
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target

